# dashboard/ghostgate-dashboard.service
[Unit]
Description=GhostGate Dashboard
After=network-online.target dnsmasq.service
Wants=network-online.target

[Service]
Type=simple
User=ggserver
WorkingDirectory=/opt/ghostgate/dashboard
Environment=GHOSTGATE_INSIDE_IF=ens37
Environment=GHOSTGATE_OUTSIDE_IF=ens33
Environment=GHOSTGATE_CLIENT_IP=10.10.10.2
Environment=GHOSTGATE_ROUTER_DNS=10.10.10.1
ExecStart=/usr/bin/python3 /opt/ghostgate/dashboard/app.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
